# target/linux/ipq40xx/base-files/etc/init.d/deco_m4r_optimize
#!/bin/sh /etc/rc.common

START=99
STOP=10

USE_PROCD=1
PROG=/usr/bin/deco_m4r_optimize

start_service() {
  local model=$(cat /tmp/sysinfo/model 2>/dev/null)
  [ "$model" != "TP-Link Deco M4R v3" ] && return 0

  echo "Starting Deco M4 V3 network optimizations..."

  setup_cpu_optimization

  setup_network_optimization

  setup_wireless_optimization

  setup_ipv6_optimization

  echo "Deco M4 V3 optimizations applied successfully"
}

setup_cpu_optimization() {
  echo "performance" >/sys/devices/system/cpu/cpu0/cpufreq/scaling_governor 2>/dev/null
  echo "performance" >/sys/devices/system/cpu/cpu1/cpufreq/scaling_governor 2>/dev/null
  echo "performance" >/sys/devices/system/cpu/cpu2/cpufreq/scaling_governor 2>/dev/null
  echo "performance" >/sys/devices/system/cpu/cpu3/cpufreq/scaling_governor 2>/dev/null

  echo 1 >/proc/irq/27/smp_affinity 2>/dev/null  # eth0
  echo 2 >/proc/irq/28/smp_affinity 2>/dev/null  # eth1
  echo 4 >/proc/irq/142/smp_affinity 2>/dev/null # wifi0
  echo 4 >/proc/irq/143/smp_affinity 2>/dev/null # wifi1

  echo 7 >/sys/class/net/eth0/queues/rx-0/rps_cpus 2>/dev/null
  echo 7 >/sys/class/net/eth1/queues/rx-0/rps_cpus 2>/dev/null
}

setup_network_optimization() {
  sysctl -w net.core.rmem_default=262144
  sysctl -w net.core.rmem_max=16777216
  sysctl -w net.core.wmem_default=262144
  sysctl -w net.core.wmem_max=16777216
  sysctl -w net.core.netdev_max_backlog=5000
  sysctl -w net.core.netdev_budget=300

  sysctl -w net.ipv4.tcp_rmem="4096 87380 16777216"
  sysctl -w net.ipv4.tcp_wmem="4096 65536 16777216"
  sysctl -w net.ipv4.tcp_congestion_control=bbr
  sysctl -w net.ipv4.tcp_slow_start_after_idle=0

  sysctl -w net.ipv4.conf.all.forwarding=1
  sysctl -w net.ipv4.conf.all.accept_redirects=0
  sysctl -w net.ipv4.conf.all.send_redirects=0

  echo 0 >/sys/class/net/br-lan/bridge/multicast_snooping 2>/dev/null
  echo 1 >/sys/class/net/br-lan/bridge/stp_state 2>/dev/null
}

setup_wireless_optimization() {
  [ -d /sys/kernel/debug/ieee80211/phy0 ] && {
    echo 0 >/sys/kernel/debug/ieee80211/phy0/ath10k/coex_bt_rssi_limit 2>/dev/null
    echo 1 >/sys/kernel/debug/ieee80211/phy0/ath10k/mesh_enabled 2>/dev/null
  }

  [ -d /sys/kernel/debug/ieee80211/phy1 ] && {
    echo 0 >/sys/kernel/debug/ieee80211/phy1/ath10k/coex_bt_rssi_limit 2>/dev/null
    echo 1 >/sys/kernel/debug/ieee80211/phy1/ath10k/mesh_enabled 2>/dev/null
  }

  iw dev wlan0 set power_save off 2>/dev/null
  iw dev wlan1 set power_save off 2>/dev/null
}

setup_ipv6_optimization() {
  sysctl -w net.ipv6.conf.all.forwarding=1
  sysctl -w net.ipv6.conf.all.accept_ra=1
  sysctl -w net.ipv6.conf.all.accept_redirects=0
  sysctl -w net.ipv6.conf.default.forwarding=1
  sysctl -w net.ipv6.conf.default.accept_ra=1

  sysctl -w net.ipv6.neigh.default.gc_thresh1=128
  sysctl -w net.ipv6.neigh.default.gc_thresh2=512
  sysctl -w net.ipv6.neigh.default.gc_thresh3=1024

  sysctl -w net.ipv6.conf.all.use_tempaddr=2
  sysctl -w net.ipv6.conf.default.use_tempaddr=2
}

stop_service() {
  local model=$(cat /tmp/sysinfo/model 2>/dev/null)
  [ "$model" != "TP-Link Deco M4 V3" ] && return 0

  echo "Stopping Deco M4 V3 optimizations..."
}

restart() {
  stop
  start
}
